# bgapi-k3s/bgapi-watch.service
# bgapi-watch.service — one-shot service that runs watch-deploy.sh.
#
# This is triggered by bgapi-watch.timer, not run manually.
# The script checks for new commits; if none, it exits immediately.
[Unit]
Description=BGAPI — poll git remote and redeploy on new commits
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/BGAPI
ExecStart=/home/pi/bgapi-k3s/watch-deploy.sh

# Log to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bgapi-watch

# Timeout: builds can take 15+ min on a Pi
TimeoutStartSec=1800
