# bgapi-k3s/bgapi-webhook.service
# bgapi-webhook.service — real-time CD receiver.
#
# Runs the `webhook` binary (sudo apt install -y webhook) as a persistent
# service listening on localhost:9000. Azure DevOps Service Hooks calls it
# through the Cloudflare Tunnel on every push; a matching HMAC signature
# triggers /home/pi/bgapi-k3s/watch-deploy.sh (build → import → rollout).
#
# This is the push-based alternative to bgapi-watch.timer (5-min polling).
# Use one or the other — not both.
[Unit]
Description=BGAPI — webhook receiver for real-time CD
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/BGAPI
ExecStart=/usr/bin/webhook -hooks /etc/webhook.conf -port 9000 -verbose
Restart=on-failure
RestartSec=5

# Log to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bgapi-webhook

[Install]
WantedBy=multi-user.target
